# docker/supervisord/kafka.ini
[program:kafka]
command=start-kafka.sh /config/server.properties
autostart=true
autorestart=true