# docker/supervisord/zookeeper.ini
[program:zookeeper]
command=/srv/kafka/bin/zookeeper-server-start.sh /srv/kafka/config/zookeeper.properties
autostart=true
autorestart=true